# conf/mariadb.conf
[program:mariadb]
command=/usr/sbin/mysqld
stdout_events_enabled=true
stderr_events_enabled=true
